# cilium-router/rootfs/etc/services.d/cilium-agent/finish
#!/usr/bin/with-contenv bashio
bashio::log.warning "cilium-agent stopped (exit code: ${1})"
