# cilium-router/rootfs/etc/services.d/node-heartbeat/run
#!/usr/bin/with-contenv bashio

export KUBECONFIG=/etc/cilium/kubeconfig
NODE_NAME=$(bashio::config 'node_name')

bashio::log.info "Starting node heartbeat for '${NODE_NAME}'..."

while true; do
    TIMESTAMP=$(date -u +"%Y-%m-%dT%H:%M:%SZ")

    # Update node lease (keeps node "Ready")
    kubectl patch lease "${NODE_NAME}" -n kube-node-lease \
        --type=merge \
        -p "{\"spec\":{\"renewTime\":\"${TIMESTAMP}\"}}" \
        2>/dev/null || \
    kubectl create -f - <<EOF 2>/dev/null || true
apiVersion: coordination.k8s.io/v1
kind: Lease
metadata:
  name: ${NODE_NAME}
  namespace: kube-node-lease
spec:
  holderIdentity: ${NODE_NAME}
  leaseDurationSeconds: 40
  renewTime: "${TIMESTAMP}"
EOF

    # Update node Ready condition
    kubectl patch node "${NODE_NAME}" --type=merge --subresource=status \
        -p "{\"status\":{\"conditions\":[{\"type\":\"Ready\",\"status\":\"True\",\"lastHeartbeatTime\":\"${TIMESTAMP}\",\"lastTransitionTime\":\"${TIMESTAMP}\",\"reason\":\"CiliumRouterRunning\",\"message\":\"Cilium agent running\"}]}}" \
        2>/dev/null || true

    sleep 10
done
